# systemd/beast-backend.service
[Unit]
Description=Boost.Beast API backend (owns and supervises nginx)
After=network.target
# The backend starts/stops/reloads nginx itself via webengine::NginxController.
# The distro nginx.service must NOT also manage nginx, or you get a double-start
# (and as root, breaking the same-user control model). Mask it — see README.md.
# Conflicts ensures the two never run at once if nginx.service is still present.
Conflicts=nginx.service

[Service]
Type=simple
User=www-data
Group=www-data
# Clean stale runtime state so bind()/pidfile checks start fresh (mirrors
# entrypoint.sh): the socket, the nginx pidfile, and the nginx temp dir.
ExecStartPre=/bin/rm -f /tmp/backend.sock /tmp/nginx.pid
ExecStartPre=/bin/mkdir -p /tmp/nginx
ExecStart=/usr/local/bin/backend
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=multi-user.target
